# templates/mountfile.automount.j2
[Unit]
Description=automount {{ item.share }} share

[Automount]
Where=/{{ mnt_point }}/{{ item.mount }}
TimeoutIdleSec={{ timeout_in_sec }}

[Install]
WantedBy={{ run_lvl }}
